# bot.conf
[program:neuro_bot_sk]
directory=/root/skbot
command = python3 bot.py

autostart=true
autorestart=true
stderr_logfile=/root/crypto_airdrop_referral_bot/logfile.log