# bot.conf
[program:neuro_bot_sk]
directory=/root/skbot
command = python3.10 bot.py

autostart=true
autorestart=true
stderr_logfile=/root/skbot/logfile.log